# com.zutubi.pulse.servercore/src/bin/init.sh.in
#! /bin/sh

# Pulse script for SysV-style init.  Requires a jsvc binary from the Apache
# Commons Daemon project (http://commons.apache.org/proper/commons-daemon/).
# Most Linux distributions provide a jsvc or apache-commons-daemon-jsvc
# package, otherwise you can build binaries with a GNU toolchain.

# You *MUST* at least define PULSE_HOME, pointing to the pulse installation
# directory.  You may also want to define other environment variables.
#
# The preferred way of modifying the Pulse environment is by creating one of
# the below files and adding variables there. (These will persist across
# upgrades, unlike direct changes to this script.)
#-----------------------------------------------------------------------------

# Pick up pulse config if it exists
if [ -f /etc/zutubi-__APP_NAME__ ]
then
    . /etc/zutubi-__APP_NAME__
fi

# Pick up pulse sysconfig (RedHat-style)
if [ -f /etc/sysconfig/zutubi-__APP_NAME__ ]
then
    . /etc/sysconfig/zutubi-__APP_NAME__
fi

# If you prefer to define PULSE_HOME directly, uncomment and edit below.
#PULSE_HOME=/usr/local/pulse

# Uncomment to run Pulse as the specified user.  Make sure the user has
# read/write permission to the Pulse installation and data directories.
#PULSE_USER=pulse

# LSB-standard init script info usable by tools such as chkconfig and
# update-rc.d to install symlinks in relevant runlevel directories.
### BEGIN INIT INFO
# Provides: zutubi-__APP_NAME__
# Required-Start: $local_fs $network
# Required-Stop: $local_fs $network
# Default-Start:  2 3 4 5
# Default-Stop: 0 1 6
# Description: __SERVICE_DESCRIPTION__
### END INIT INFO

# Do not modify anything beyond this point
#-----------------------------------------------------------------------------
APP_NAME="__APP_NAME__"
APP_LONG_NAME="__SERVICE_NAME__"

if [ -z "$PULSE_HOME" ]
then
    echo "PULSE_HOME is not defined."
    exit 1
fi

if [ ! -d "$PULSE_HOME" ]
then
    echo "PULSE_HOME ($PULSE_HOME) is not a valid directory."
    exit 1
fi

if [ -z "$PULSE_JSVC" ]
then
    PULSE_JSVC=`which jsvc`
elif [ ! -f "$PULSE_JSVC" ]
then
    echo "PULSE_JSVC ($PULSE_JSVC) does not exist."
    exit 1
fi

if [ -z "$PULSE_JSVC" ]
then
    echo "PULSE_JSVC not defined and jsvc not found in the PATH."
    echo "Please install jsvc and/or add it to the system PATH."
    exit 1
fi

if [ -z "$PULSE_USER" ]
then
    USER_ARGS=
else
    USER_ARGS="-user $PULSE_USER"
fi

if [ -z "$PULSE_PIDFILE" ]
then
    PULSE_PIDFILE="$PULSE_HOME/bin/$APP_NAME.pid"
fi

if [ -z "$PULSE_OUTFILE" ]
then
    PULSE_OUTFILE="$PULSE_HOME/logs/$APP_NAME.out"
fi

if [ -z "$PULSE_OPTS" ]
then
    PULSE_OPTS="-Xmx1024m -XX:MaxPermSize=128m"
fi

PULSE_CLASSPATH="$PULSE_HOME/lib/boot.jar:$PULSE_HOME/lib/commons-daemon.jar"

start() {
    echo "Starting $APP_LONG_NAME..."
    "$PULSE_JSVC" \
        $USER_ARGS \
        -pidfile "$PULSE_PIDFILE" \
        -outfile "$PULSE_OUTFILE" \
        -errfile '&1' \
        -server \
        -Dpulse.home="$PULSE_HOME" \
        -Djava.awt.headless=true \
        -Djava.util.logging.config.class=com.zutubi.pulse.logging.ConsoleConfig \
        $PULSE_OPTS \
        -cp "$PULSE_CLASSPATH" \
        com.zutubi.pulse.command.PulseDaemon start stopservice
}

stop() {
    echo "Stopping $APP_LONG_NAME..."
    "$PULSE_JSVC" \
        -pidfile "$PULSE_PIDFILE" \
        -Dpulse.home="$PULSE_HOME" \
        -Djava.awt.headless=true \
        -stop \
        -cp "$PULSE_CLASSPATH" \
        com.zutubi.pulse.command.PulseDaemon
}

case "$1" in
    start)
        start
        ;;

    stop)
        stop
        ;;
    
    restart)
        stop
        start
        ;;
    *)
        echo "Usage: $0 { start |stop | restart }"
        exit 1;;
esac

exit 0
